# apps/edge-gateway/kitnets-gateway.service
[Unit]
Description=Kitnets Smart Meter Gateway
After=network.target

[Service]
Type=simple
User=root
# Change User to appropriate non-root user if needed (e.g. ubuntu or pi)
# If using valid port < 1024, might need CAP_NET_BIND_SERVICE or root.
# App uses port 3000 by default.

WorkingDirectory=/opt/kitnets-gateway
ExecStart=/usr/bin/npm start
# OR direct node: ExecStart=/usr/bin/node dist/index.js

Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production

# Load env vars from .env file
EnvironmentFile=/opt/kitnets-gateway/.env

[Install]
WantedBy=multi-user.target
